# mio_bot_common_ports.service
[Unit]
Description=Network Common Ports
After=multi-user.target

[Service]
Type=idle
ExecStart=/usr/bin/python3 /home/pi/bot_common_ports/bot_common_ports.py >> /home/pi/bot_common_ports/bot.log 2>&1
User=pi

[Install]
WantedBy=multi-user.target
